# target/supervisor/conf.d/webserver.conf
[program:webserver]
command=python3 /usr/local/bin/healthcheck.py
autostart=true
autorestart=true
stderr_logfile=/var/log/webserver.err.log
stdout_logfile=/var/log/webserver.out.log